# csm/conf/etc/cron.daily/es_logrotate.cron
#!/bin/bash

FILE=/usr/bin/csm_cleanup
if [ -f "$FILE" ]; then
    $FILE auditlogs -i csmauditlog -d 7
fi
