# csm/conf/etc/cron.daily/es_logrotate.cron
#!/bin/bash

FILE=/opt/seagate/cortx/csm/conf/es_data_cleanup.py
if [ -f "$FILE" ]; then
    python3 $FILE -i csmauditlog -d 7
fi
